# ops/systemd/nanobot-llm-warmup.service
[Unit]
Description=Warm nanobot LLM cache after container startup
After=podman-nanobot-cage.service
Requires=podman-nanobot-cage.service

[Service]
Type=oneshot
TimeoutStartSec=80
ExecStart=/bin/sh -lc 'sleep 30; /usr/bin/podman exec nanobot-cage timeout 35 python3 -m nanobot.agent.warmup --limit 1 --timeout 25 --max-tokens 4 || echo "LLM warmup best-effort attempt failed or timed out"; exit 0'
