# ops/systemd/obp-rs.service
[Unit]
Description=OBP Failover Rust Service
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/root/obp-rs
Environment=OBP_HOST=127.0.0.1
Environment=OBP_PORT=8000
ExecStart=/usr/local/bin/obp-rs
Restart=always
RestartSec=3
User=root

[Install]
WantedBy=multi-user.target
