# ops/systemd/trend-sidecar-rs.service
[Unit]
Description=Trend Radar Rust Sidecar
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
Environment=TREND_SIDECAR_PORT=8095
Environment=TREND_REFRESH_SECS=1800
Environment=TREND_STATE_FILE=/root/.nanobot/data/trend-sidecar/state.json
ExecStart=/usr/local/bin/trend-sidecar-rs
Restart=always
RestartSec=5
User=root
WorkingDirectory=/root/nanobot-ops/sources/trend-sidecar-rs

[Install]
WantedBy=multi-user.target
