# ops/systemd/drop-ins/lof-sidecar.service.d/20-obp-auth.conf
[Service]
EnvironmentFile=-/root/.nanobot/secrets/obp-proxy.env
