# ops/systemd/drop-ins/qq-sidecar-rs.service.d/20-podman-bridge.conf
[Service]
Environment=QQ_SIDECAR_HOST=172.17.0.1
