# ops/systemd/drop-ins/podman-nanobot-cage.service.d/20-llm-warmup.conf
[Unit]
Wants=nanobot-llm-warmup.service
